# workspace/crontab/laradock
* * * * * laradock /usr/bin/php /var/www/laragis/artisan schedule:run >> /dev/null 2>&1
